# board/nano_linux/rootfs-overlay/etc/profile.d/start_gui.sh
#!/bin/sh
# Если мы в первой консоли tty1 и X еще не запущен
if [ -z "$DISPLAY" ] && [ "$(tty)" = "/dev/tty1" ]; then
    # Копируем конфиг, если его нет (хотя лучше сделать это через симлинк)
    # И запускаем X-сервер
    exec startx
fi
